# debian/extra/user-units/graphical-session-pre.target
[Unit]
Description=Session services which should run before the graphical session is brought up
StopWhenUnneeded=yes
