# debian/libnss-resolve.postinst
#!/bin/sh

set -e

if [ "$1" = configure ] && [ -z "$2" ]; then
    echo "First installation detected..."
    # Enable resolved
    systemctl enable systemd-resolved.service
    if [ -d /run/systemd/system ]; then
        deb-systemd-invoke start systemd-resolved.service || true
    fi
fi

#DEBHELPER#
