# debian/libnss-resolve.postrm
#!/bin/sh

set -e

if [ "$1" = remove ]; then
    systemctl disable systemd-resolved.service
    if [ -d /run/systemd/system ]; then
        deb-systemd-invoke stop systemd-resolved.service || true
    fi
fi

#DEBHELPER#
