# debian/extra/units/systemd-logind.service.d/nice.conf
[Service]
# Set Nice=-1 to doge watchdog on soft lockups (LP: #1696970).
Nice=-1
